# dotfiles/.config/systemd/user/redshift.service
[Unit]
Description=Redshift color temperature adjustment

[Service]
ExecStart=/usr/bin/redshift
Restart=always

[Install]
WantedBy=default.target
